# pipeline/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:qbittorrent]
command=qbittorrent-nox --webui-port=8080 --confirm-legal-notice
autostart=true
autorestart=true
stdout_logfile=/tmp/qbit.log
stderr_logfile=/tmp/qbit.err

[program:worker]
command=bash -c "exec uvicorn worker:app --host 0.0.0.0 --port $${PORT:-8090}"
directory=/app/pipeline
autostart=true
autorestart=true
stdout_logfile=/tmp/worker.log
stderr_logfile=/tmp/worker.err
